# joplin-anki-sync.cron
#/etc/cron.d/joplin-anki-sync.cron

PATH=/bin:/usr/bin:/usr/sbin:/usr/local/bin
*/10 * * * * andrew joplin-anki-sync >/dev/null 2>&1
